# xray-subscription.service
[Unit]
Description=Xray Subscription Server
Documentation=https://github.com/XTLS/Xray-core
After=network.target xray.service

[Service]
Type=simple
# Run as xray user so Raven and Xray share ownership of config_dir files.
# When api_user_inbound_tag is set, Raven writes to config_dir; Xray must read those files.
User=xray
Group=xray
ExecStart=/usr/local/bin/xray-subscription -config /etc/xray-subscription/config.json
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/xray-subscription /etc/xray/config.d
ReadOnlyPaths=/etc/xray-subscription
PrivateTmp=yes
ProtectHome=yes

[Install]
WantedBy=multi-user.target
